# packaging/buxton2.socket
[Unit]
Description=Buxton Configuration Service

[Socket]
ListenStream=/run/buxton2-0
SocketMode=0777
SmackLabelIPIn=*
SmackLabelIPOut=@

[Install]
WantedBy=sockets.target
